# deamon-service/msocksd.service
[Service]
ExecStart=/usr/bin/msocks
Type = simple
